# files/riak.conf
description "Start a Riak node"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]
respawn
respawn limit 10 5

console output

exec /usr/sbin/riak start